# texdocker
#!/sbin/openrc-run

depend() {
    need docker
}

stop() {
    docker stop tex_container
}

# command=/home/$USER/projects/texdocker/run-docker.sh
