# texdocker
#!/sbin/openrc-run

depend() {
    need docker
}

# command=/home/$USER/projects/texdocker/run-docker.sh
