# texdocker
#!/sbin/openrc-run

start() {
    echo "starting texdocker.."
    su $command_user -c "$command $command_args"
}

stop() {
    echo "stopping texdocker..."
    podman stop tex_container
    echo "stopped texdocker"
}

