# texman
#!/sbin/openrc-run

depend() {
    need cgroups
}

start() {
    echo "starting texman..."
    # important to start as $USER / rootless
    su $command_user -c "$command $command_args"
}

stop() {
    echo "stopping texman..."
    # important to stop as $USER / rootless
    su $command_user -c "podman stop tex_container"
    echo "stopped texman"
}

