# texman
#!/sbin/openrc-run

depend() {
    need cgroups
}

start() {
    echo "STARTING TEXMAN"
    max_attempts=30
    interval=.2
    attempt=1
    while [ $attempt -le $max_attempts ]; do
        # call run-podman.sh
        # important to start as $USER / rootless
        su $command_user -c "$command $command_args"
        if [ $? -eq 0 ]; then
            exit 0
        fi
        echo "Attempt $attempt failed. Retrying..."
        sleep $interval
        attempt=$((attempt + 1))
    done
}

stop() {
    echo "STOPPING TEXMAN"
    # important to stop as $USER / rootless
    su $command_user -c "podman stop tex_container"
    echo "TEXMAN STOPPED"
}